# deploy/systemd/openclaw-control-plane.service
[Unit]
Description=OpenClaw AutoDeploy Control Plane API
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
User=openclaw
Group=openclaw
WorkingDirectory=/opt/openclaw-autodeploy
EnvironmentFile=/etc/openclaw-autodeploy/control-plane.env
ExecStart=/opt/openclaw-autodeploy/bin/control-plane-api --config /etc/openclaw-autodeploy/config.yaml
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
